# site-cookbooks/brewbit/templates/default/etc/cron.daily/backup.erb
#!/bin/sh
su - deploy -c "cd /var/www/<%= node[:brewbit][:hostname] %>/current/backup/ && bundle install --path vendor/bundle && RAILS_ENV=<%= node.chef_environment %> bundle exec backup perform --root-path . --trigger db"
